# debian/suspend-prevent.service
[Unit]
Description=Prevent suspend on lid close when on AC

[Service]
Type=simple
ExecStart=/usr/bin/suspend-prevent --forever
